# packaging/debian/prerm
#!/bin/bash
set -e

case "$1" in
    remove|purge)
        # Stop services
        if [ -d /run/systemd/system ]; then
            systemctl stop rune.service >/dev/null 2>&1 || true
            systemctl stop rune.socket >/dev/null 2>&1 || true
            systemctl stop rune-registry.service >/dev/null 2>&1 || true
            systemctl disable rune.socket >/dev/null 2>&1 || true
            systemctl disable rune.service >/dev/null 2>&1 || true
            systemctl disable rune-registry.service >/dev/null 2>&1 || true
        fi

        # Remove socket file
        rm -f /var/run/rune.sock
        ;;
esac

#DEBHELPER#

exit 0
